# compose/app_consumer/supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
pidfile = /tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = /var/run/supervisor.sock

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:celery-worker-consumer]
directory=/app
command=/usr/local/bin/celery --app=config.celery:celery worker --loglevel=INFO -f /var/log/celery.log --autoscale=10,0
autostart=true
autorestart=true
numprocs=1
